# docker/supervisord.app.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:python-sidecar]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8081 --workers 1
directory=/app/python-sidecar
autostart=true
autorestart=true
stdout_logfile=/var/log/python-sidecar.log
stderr_logfile=/var/log/python-sidecar_err.log
startsecs=10
startretries=5
priority=20

[program:go-backend]
command=/app/docker/wait-for-port.sh 127.0.0.1 8081 python-sidecar /app/server
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/go-backend.log
stderr_logfile=/var/log/go-backend_err.log
startsecs=5
startretries=10
priority=30
environment=GIN_MODE="release"
