# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# ── PostgreSQL 16 ────────────────────────────────────────
[program:postgresql]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql_err.log
priority=10
startsecs=3
startretries=5

# ── Redis 7 ─────────────────────────────────────────────
[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --port 6379 --loglevel warning --save 900 1 --save 300 10 --dbfilename dump.rdb --dir /var/lib/redis
autostart=true
autorestart=true
stdout_logfile=/var/log/redis.log
stderr_logfile=/var/log/redis_err.log
priority=20
startsecs=3
startretries=5

# ── Neo4j 5 CE ───────────────────────────────────────────
# Neo4j is optional: if it fails, the sidecar continues without graph support.
[program:neo4j]
command=/opt/neo4j/bin/neo4j console
user=neo4j
autostart=true
autorestart=unexpected
stdout_logfile=/var/log/neo4j.log
stderr_logfile=/var/log/neo4j_err.log
priority=15
startsecs=30
startretries=3
environment=NEO4J_HOME="/opt/neo4j",NEO4J_AUTH="neo4j/%(ENV_NEO4J_PASSWORD)s",NEO4J_server_bolt_listen__address="0.0.0.0:7687",NEO4J_server_http_listen__address="0.0.0.0:7474",NEO4J_server_memory_heap_initial__size="256m",NEO4J_server_memory_heap_max__size="512m",NEO4J_server_memory_pagecache_size="128m"

# ── Qdrant 1.12 ──────────────────────────────────────────
[program:qdrant]
command=/usr/local/bin/qdrant --storage-path /var/lib/qdrant/storage
directory=/var/lib/qdrant
autostart=true
autorestart=unexpected
stdout_logfile=/var/log/qdrant.log
stderr_logfile=/var/log/qdrant_err.log
priority=25
startsecs=5
startretries=5
environment=QDRANT__SERVICE__HTTP_PORT="6333",QDRANT__SERVICE__GRPC_PORT="6334",QDRANT__LOG_LEVEL="WARN"

# ── Python sidecar (LangGraph agent + FastAPI) ────────────
[program:python-sidecar]
command=/app/docker/wait-for-pg.sh python -m uvicorn main:app --host 0.0.0.0 --port 8081 --workers 1
directory=/app/python-sidecar
autostart=true
autorestart=true
stdout_logfile=/var/log/python-sidecar.log
stderr_logfile=/var/log/python-sidecar_err.log
startsecs=10
startretries=5
priority=30
environment=REDIS_URL="redis://127.0.0.1:6379",NEO4J_URI="bolt://127.0.0.1:7687",QDRANT_URL="http://127.0.0.1:6333"

# ── Go backend (Gin API gateway) ─────────────────────────
[program:go-backend]
command=/app/docker/wait-for-pg.sh /app/server
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/go-backend.log
stderr_logfile=/var/log/go-backend_err.log
startsecs=5
startretries=10
priority=40
environment=GIN_MODE="release"

