# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql_err.log
priority=10

[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --port 6379 --loglevel warning
autostart=true
autorestart=true
stdout_logfile=/var/log/redis.log
stderr_logfile=/var/log/redis_err.log
priority=20

[program:python-sidecar]
command=python -m uvicorn main:app --host 0.0.0.0 --port 8081
directory=/app/python-sidecar
autostart=true
autorestart=true
stdout_logfile=/var/log/python-sidecar.log
stderr_logfile=/var/log/python-sidecar_err.log
startsecs=5
startretries=3
priority=30
environment=DATABASE_URL="postgresql://%(ENV_DB_USER)s:%(ENV_DB_PASSWORD)s@%(ENV_DB_HOST)s:%(ENV_DB_PORT)s/%(ENV_DB_NAME)s"

[program:go-backend]
command=/app/server
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/go-backend.log
stderr_logfile=/var/log/go-backend_err.log
startsecs=5
startretries=3
priority=40
environment=GIN_MODE="release"
